# mybot.conf
[program:elonbot-org]
directory=/home/ubuntu/bots/elonbot-org
command=python3 app.py

autostart=true
autorestart=true
environment=HOME="/home/ubuntu",USER="ubuntu"

stderr_logfile=/home/ubuntu/bots/elonbot-org/logfile_err.log
stdout_logfile=/home/ubuntu/bots/elonbot-org/logfile.log